# launchd/com.multinexus.mac-codex.agentd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.multinexus.mac-codex.agentd</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yinxin/projects/multinexus/.venv/bin/python</string>
        <string>-m</string>
        <string>multinexus.agentd</string>
        <string>--config</string>
        <string>/Users/yinxin/projects/multinexus/agents.toml</string>
        <string>--agent</string>
        <string>mac-codex</string>
        <string>--poll-interval</string>
        <string>2</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yinxin/projects/multinexus</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>10</integer>
    <key>StandardOutPath</key>
    <string>/Users/yinxin/projects/multinexus/logs/mac-codex.agentd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yinxin/projects/multinexus/logs/mac-codex.agentd.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yinxin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Users/yinxin/.local/bin:/Users/yinxin/.nvm/versions/node/v24.14.1/bin:/Users/yinxin/.opencode/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
