# launchd/com.discord-nexus.mac-omp.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.discord-nexus.mac-omp</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/yinxin/projects/discord-nexus/.venv/bin/python</string>
        <string>/Users/yinxin/projects/discord-nexus/nexus.py</string>
        <string>--agent</string>
        <string>mac-omp</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/yinxin/projects/discord-nexus</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>
    <key>StandardOutPath</key>
    <string>/Users/yinxin/projects/discord-nexus/logs/mac-omp.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/yinxin/projects/discord-nexus/logs/mac-omp.err.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/yinxin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
        <key>PATH</key>
        <string>/Users/yinxin/.local/bin:/Users/yinxin/.bun/bin:/Users/yinxin/.nvm/versions/node/v24.14.1/bin:/Users/yinxin/.opencode/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>
</dict>
</plist>
